# package/lean/luci-app-sysmonitor/root/etc/hotplug.d/iface/90-sysmonitor
[ $INTERFACE = wan6 ] && touch /tmp/wan6.sign
[ $INTERFACE = wan ] && touch /tmp/wan.sign

